# configs/supervisor.d/dkim.ini
[program:dkim]
command=/usr/bin/opendkim -x /etc/opendkim/opendkim.conf -f
redirect_stderr=true
stdout_logfile=/var/log/supervisor/dkim.log
