# configs/supervisor.d/dmarc.ini
[program:dmarc]
command=/usr/bin/opendmarc -c /etc/opendmarc/opendmarc.conf -f
redirect_stderr=true
stdout_logfile=/var/log/supervisor/dmarc.log
