# configs/supervisor.d/dovecot.ini
[program:dovecot]
command=/usr/bin/dovecot -F
redirect_stderr=true
stdout_logfile=/data/log/supervisor/dovecot.log
