# configs/supervisor.d/mysql.ini
[program:mysql]
user=mysql
command=/usr/sbin/mysqld --datadir /data/mysql
redirect_stderr=true
stdout_logfile=/var/log/supervisor/mysql.log
