# configs/supervisor.d/nginx.ini
[program:nginx]
command=/usr/bin/nginx -g 'daemon off;'
redirect_stderr=true
stdout_logfile=/data/log/supervisor/nginx.log
