# configs/supervisor.d/php.ini
[program:php]
command=/usr/bin/php-fpm --nodaemonize --fpm-config /etc/php/php-fpm.conf
redirect_stderr=true
stdout_logfile=/var/log/supervisor/php.log
