# configs/supervisor.d/postfix.ini
[program:postfix]
command=/scripts/postfix.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/postfix.log
