# configs/supervisor.d/spam.ini
[program:spam]
command=/usr/bin/vendor_perl/spamd -x -u spamd -g spamd
redirect_stderr=true
stdout_logfile=/var/log/supervisor/spam.log
