# configs/supervisor.d/srs.ini
[program:srs]
command=/usr/sbin/postsrsd -s/data/postsrsd/postsrsd.secret -c/usr/lib/postsrsd -upostsrsd -e
environment=SRS_DOMAIN="__DOMAIN__",SRS_EXCLUDE_DOMAINS="__DOMAINS__"
redirect_stderr=true
stdout_logfile=/data/log/supervisor/srs.log
