# configs/supervisor.d/syslog.ini
[program:syslog]
command=/usr/bin/rsyslogd -n -i /run/rsyslogd.pid
redirect_stderr=true
stdout_logfile=/data/log/supervisor/syslog.log
