# configs/cron.daily/roundcube_clean.sh
#!/bin/bash

/usr/share/webapps/roundcubemail/bin/cleandb.sh
