# configs/cron.daily/spam_update.sh
#!/bin/bash
/scripts/update_sa.sh learn

