# configs/systemd/system/roundcube_clean.timer
[Unit]
Description=RoundCube CleanUp Cron Job

[Timer]
OnCalendar=daily
Persistent=true     
 
[Install]
WantedBy=timers.target
